# src/ansible/roles/boxdev/files/mount-boxdev.conf
description "Mount boxdev drives"

start on filesystem
stop on shutdown

exec /sbin/mount-boxdev.sh >> /var/log/boxdev.log 2>&1

